# kiosk/pulse-board.service
# systemd user service to auto-start the kiosk on boot (Raspberry Pi OS with a
# desktop session / X11).
#
# Don't install this by hand — kiosk/install.sh fetches it, writes your
# CloudFront URL into the Environment line below, and enables the service:
#
#   curl -fsSL https://raw.githubusercontent.com/nicholastanderson/b9-owner-dashboard/main/kiosk/install.sh \
#     | bash -s -- https://YOUR-DIST.cloudfront.net
#
# To change the URL later:
#   systemctl --user edit --full pulse-board.service
#   systemctl --user restart pulse-board.service
#
[Unit]
Description=Back Nine Pulse Board kiosk (Chromium)
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
Environment=DISPLAY=:0
# Rewritten by kiosk/install.sh with your deployed CloudFront URL:
Environment=PULSE_BOARD_URL=https://REPLACE_ME.cloudfront.net
ExecStart=%h/pulse-board-kiosk.sh
Restart=always
RestartSec=5

[Install]
WantedBy=graphical-session.target
